# warnquota-ng.crontab
4,9,14,19,24,29,34,39,44,49,54,59 * * * *	root	WQNG_BINARYPATH/warnquota-ng -q
*/5 * * * *	WQNG_USERNAME	WQNG_BINARYPATH/warnquota-ng -s
# to force UTF-8 locale use this line:
#4,9,14,19,24,29,34,39,44,49,54,59 * * * *	root	env LC_ALL=en_US.UTF-8 /usr/sbin/warnquota-ng -q
#*/5 * * * *	WQNG_USERNAME	env LC_ALL=en_US.UTF-8 WQNG_BINARYPATH/warnquota-ng -s
# Regen fs configs from LDAP
#11 4 * * *	WQNG_USERNAME	WQNG_BINARYPATH/warnquota-ng -L
